# ansible/etc/profile.d/sqoop.sh
export SQOOP_HOME=/opt/sqoop
export PATH=$PATH:$SQOOP_HOME/bin
